# src/bsd_bt_daemon.rc
#!/bin/sh
#
# PROVIDE: bsd_bt_daemon
# REQUIRE: DAEMON hcsecd
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable bsd_bt_daemon:
# bsd_bt_daemon_enable="YES"
#
# Configuration variables:
# bsd_bt_daemon_user (str):	User to run the daemon as (default: root)
# bsd_bt_daemon_flags (str):	Flags to pass to the daemon (e.g., "--device ubt0hci")

. /etc/rc.subr

name="bsd_bt_daemon"
rcvar="bsd_bt_daemon_enable"

load_rc_config $name

: ${bsd_bt_daemon_enable:="NO"}
: ${bsd_bt_daemon_user:="root"}
: ${bsd_bt_daemon_flags:=""}

# Using the entry point installed by pip/setup.py
command="/usr/local/bin/bsd-bt-daemon"

# PID file handling
pidfile="/var/run/${name}.pid"

start_cmd="${name}_start"

bsd_bt_daemon_start() {
	if [ ! -x "${command}" ]; then
		echo "Error: ${command} not found or not executable."
		return 1
	fi
	
	echo "Starting ${name}..."
	# Run with daemon utility to background it and handle PID file
	/usr/sbin/daemon -u "${bsd_bt_daemon_user}" -p "${pidfile}" -f "${command}" ${bsd_bt_daemon_flags}
}

run_rc_command "$1"